# src/profile.d/system-info.csh
# Only run System-Info if this session has a terminal:
test -t 1 && System-Info
