# src/profile.d/system-info.sh
# Only run System-Info if this session has a terminal:
if [ -t 1 ] ; then
   System-Info
fi
